# balena-ads-b/planefinder/supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:pfclient]
command=/usr/bin/pfclient
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
